# backend/mealtogether-celery.service
[Unit]
Description=MealTogether Celery Worker
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/applications/meal-together/backend
Environment="PATH=/opt/applications/meal-together/backend/venv/bin"
Environment="FLASK_ENV=production"
EnvironmentFile=/opt/applications/meal-together/backend/.env
ExecStart=/opt/applications/meal-together/backend/venv/bin/celery -A celery_app.celery worker --loglevel=info
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
